# docker/grpc/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:grpc_server]
command=python3.7 /opt/app/server.py
user=root
autostart=true
